# debian/prerm
#!/bin/bash
set -e

# Stop and disable the service
systemctl stop pypltagent.service || true
systemctl disable pypltagent.service || true

echo "pypltagent service stopped and disabled."
